# custom_media_agent_2_0/deploy/systemd/alchemy-v2-sync-worker.service
[Unit]
Description=Alchemy Custom Media Agent 2.0 Resource Sync Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alchemy
Group=alchemy
WorkingDirectory=/opt/alchemy/AlchemyOS/custom_media_agent_2_0
EnvironmentFile=/etc/alchemy/alchemy-v2.env
ExecStart=/opt/alchemy/AlchemyOS/custom_media_agent_2_0/.venv/bin/python -m app.workers.resource_sync_worker --mode auto
Restart=always
RestartSec=30
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
